# magic-trackpad-battery.service
[Unit]
Description=Magic Trackpad battery monitor (Bluetooth HID)
After=bluetooth.target

[Service]
ExecStart=%h/.local/bin/magic-trackpad-battery
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
RestrictSUIDSGID=true
LockPersonality=true
PrivateTmp=true

[Install]
WantedBy=default.target
