# magic-trackpad-autoconnect.service
[Unit]
Description=Auto-connect paired Magic Trackpad
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=oneshot
Environment=MAX_ATTEMPTS=3
Environment=SLEEP_SECONDS=2
Environment=CONNECT_TIMEOUT=4
ExecStart=%h/.local/bin/magic-trackpad-connect
SuccessExitStatus=1
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
RestrictSUIDSGID=true
LockPersonality=true
PrivateTmp=true
